# django/deploy/supervisord/supervisord.conf
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true
logfile=/miracle/logs/supervisor/supervisord.log
user=miracle
childlogdir=/miracle/logs/supervisor            ; ('AUTO' child log dir, default $TEMP)

[program:uwsgi]
command=uwsgi --ini /code/django/deploy/uwsgi/miracle.ini
stdout_logfile=/miracle/logs/uwsgi.log
redirect_stderr=true
directory=/code/django
user=miracle
autostart=true
autorestart=true

[program:redis]
command=/usr/bin/redis-server
user=miracle
stdout_logfile=/miracle/logs/redis.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=3

[program:celery]
command=celery -A miracle worker -l debug
directory=/code/django
user=miracle
stdout_logfile=/miracle/logs/celery.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=6
